# deploy/systemd/clothing-store-api.service
[Unit]
Description=Clothing Store API
After=network.target postgresql.service
Wants=postgresql.service

[Service]
WorkingDirectory=/opt/clothing_store
ExecStart=/usr/bin/dotnet /opt/clothing_store_runtime/store-api/Store.Api.dll
Restart=always
RestartSec=5
User=www-data
StateDirectory=clothing-store-api
CacheDirectory=clothing-store-api
Environment=HOME=/var/lib/clothing-store-api
Environment=DOTNET_CLI_HOME=/var/lib/clothing-store-api/.dotnet
Environment=DOTNET_BUNDLE_EXTRACT_BASE_DIR=/var/cache/clothing-store-api/dotnet-bundle
Environment=DOTNET_SKIP_FIRST_TIME_EXPERIENCE=1
Environment=DOTNET_NOLOGO=1
Environment=DOTNET_CLI_TELEMETRY_OPTOUT=1
EnvironmentFile=/etc/clothing-store/environment

[Install]
WantedBy=multi-user.target
